# etc/pymap/pymap@.service
[Unit]
Description = pymap IMAP server - %I
Wants = pymap-%i.socket

[Service]
EnvironmentFile = /etc/default/pymap-%i
ExecStart = /opt/pymap/bin/pymap $PYMAP_ARGS

[Install]
WantedBy = multi-user.target
